# crontab
30 0 * * * python /dist/main.py /config.json
